# sample.supervisor-daemon.conf
[program:led-board-manager]
command=/home/pi/led-board-manager/env/bin/gunicorn Capstone.wsgi -b 0:9002  --capture-output
directory=/home/pi/led-board-manager
user=pi
stopasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/led-board-manager.stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:scoreboard]
command=python3 ./src/main.py --led-rows=32 --led-cols=64 --led-brightness=60 --led-gpio-mapping=adafruit-hat --led-slowdown-gpio=2 --updatecheck
directory=/home/pi/nhl-led-scoreboard
autostart=true
autorestart=true
stderr_logfile=/var/log/scoreboard.stderr.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
stdout_logfile=/var/log/scoreboard.stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10